# .github/systemd/mobilesync-server.service
[Unit]
Description=MobileSync Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mobilesync
ExecStart=/opt/mobilesync/mobilesync-server
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
